# templates/systemd/reboot.timer.j2
[Unit]
Description=periodic system reboot timer

[Timer]
OnBootSec={{ os_cfg_reboot | default("1month") }}
Unit=reboot.service

[Install]
WantedBy=timers.target