# molecule/default/files/reboot.timer
[Unit]
Description=periodic system reboot timer

[Timer]
OnBootSec=1month
Unit=reboot.service

[Install]
WantedBy=timers.target
